# relogin/systemd/bezant-relogin.timer
[Unit]
Description=Probe bezant-server every 5 min and re-login to IBKR if expired
Documentation=https://github.com/isaacrowntree/bezant/tree/main/relogin

[Timer]
# Fire 2 min after boot to give bezant-server / cloudflared time to settle,
# then every 5 min thereafter.
OnBootSec=2min
OnUnitActiveSec=5min
# If the Pi was off when a tick should have fired, run it as soon as we come back.
Persistent=true
Unit=bezant-relogin.service

[Install]
WantedBy=timers.target
