# relogin/systemd/bezant-relogin.service
[Unit]
Description=Re-login to IBKR Client Portal Gateway when the session expires
Documentation=https://github.com/isaacrowntree/bezant/tree/main/keepalive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/bezant-relogin
EnvironmentFile=%h/bezant-relogin/.env
ExecStart=/usr/bin/env npx tsx %h/bezant-relogin/index.ts
StandardOutput=journal
StandardError=journal
# Don't restart automatically — the timer handles cadence and the script
# self-disables after MAX_CONSECUTIVE_FAILURES.
Restart=no
# Bound the worst-case run (2-min IB Key wait + 30s slack)
TimeoutStartSec=180

[Install]
WantedBy=default.target
